# foot-server@.service.in
[Service]
ExecStart=@bindir@/foot --server=0
Environment=WAYLAND_DISPLAY=%i
NonBlocking=true
StandardInput=socket

[Unit]
Requires=%N.socket
Description=Foot terminal server mode for WAYLAND_DISPLAY=%i
Documentation=man:foot(1)

[Install]
WantedBy=wayland-instance@.target
